# dist/dante.cron
# Execute all miners each day at 23:45
0 5 * * *   root /usr/bin/ciacco

# Send notification mail once a week, after miner execution
0 5 * * 1   root cd /usr/share/dante/caronte && ./caronte 'http://localhost:9292/#/?theme=dark&palette=palette3&last=week&lang=en' root@localhost 
