# dist/virgilio.service
[Unit]
Description=Virgilio report server for Dante
After=syslog.target network-online.target systemd-tmpfiles-setup.service
Documentation=https://github.com/nethesis/dante

[Service]
EnvironmentFile=/etc/sysconfig/dante
ExecStart=/usr/bin/virgilio
User=nobody

[Install]
WantedBy=multi-user.target
